# configs/d1-watcher.service
[Unit]
Description=d1 RAG File Watcher
After=network.target

[Service]
Type=simple
# %h = user home; adjust path if d1 is not at ~/CodeP/d1
ExecStart=%h/CodeP/d1/scripts/watch.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
