# configs/d1-watcher.service
[Unit]
Description=d1 RAG File Watcher
After=network.target

[Service]
Type=simple
# Installed by scripts/install-watcher.sh with real path substituted
ExecStart=@D1_PROJECT_ROOT@/scripts/watch.sh
WorkingDirectory=@D1_PROJECT_ROOT@
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
